# mpmusic.conf
[program:mpmusic]
command=pipenv run gunicorn -w 4 manage:app
directory=/home/ubuntu/mp-music
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true